# systemd/fedora43-gaming-firstboot.service
[Unit]
Description=Fedora 43 KDE Gaming - Firstboot Provisioning
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/fedora43-gaming/scripts/firstboot/firstboot.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
